# etc/systemd/system/birdpi-wifi-onboard.service
[Unit]
Description=BirdPi WiFi Onboarding and Captive Portal
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/birdpi-wifi-onboard.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
